# templates/default/mon-agent-init.conf.erb
# Startup script for the Mon-Agent

description "Mon-Agent application"
start on runlevel [2345]

console log
respawn

setgid mon_agent
setuid mon_agent
exec /usr/bin/python /usr/lib/python2.7/dist-packages/monagent/agent.py start
